# newsboat/.local/share/systemd/user/newsboat-reloader.service
[Unit]
Description=Newsboat Reloader

[Service]
ExecStart=/usr/bin/newsboat -x reload
