# setup/crontab
*/10 * * * * pidof -x newsboat > /dev/null || newsboat -x reload -x print-unread | cut -f1 -d' ' > /tmp/newslist
